# launch-agent/com.spendpulse.app.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.spendpulse.app</string>
  <key>ProgramArguments</key>
  <array>
    <string>__SPENDPULSE_ROOT__/scripts/run-server.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>__SPENDPULSE_ROOT__</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>SPENDPULSE_ROOT</key>
    <string>__SPENDPULSE_ROOT__</string>
    <key>HOST</key>
    <string>127.0.0.1</string>
    <key>PORT</key>
    <string>8787</string>
    <key>PATH</key>
    <string>__SPENDPULSE_ROOT__/.venv/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>__SPENDPULSE_ROOT__/logs/server.out.log</string>
  <key>StandardErrorPath</key>
  <string>__SPENDPULSE_ROOT__/logs/server.err.log</string>
</dict>
</plist>
